# scripts/moist-bot.service
[Unit]
Description=Moistbot Fleabot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/moist-bot
Environment=HOME=/home/ec2-user
Environment=PATH=/home/ec2-user/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/ec2-user/.local/bin/uv run app
Restart=always
RestartSec=10
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=moist-bot

[Install]
WantedBy=multi-user.target
